# resources/templates/spot/assets/lib.systemd.system.spot.service
[Unit]
Description=Start fake spot calculation

[Service]
Type=forking
ExecStart=/usr/bin/bash -x /home/ec2-user/send_metrics
